# crontab
*/10 * * * * xvfb-run calibredb add /import -r --with-library /books && rm /import*
# Don't remove the empty line at the end of this file. It is required to run the cron job
